# drivereset/drivereset.service
[Unit]
Description=Reset PiRoCon

[Service]
Type=oneshot
ExecStart=/usr/local/bin/drivereset/drivereset.sh

[Install]
WantedBy=multi-user.target